# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#SD
mount -t vfat /dev/mmcblk0p1 /mnt/mmc

#start sound
aplay /etc/start.wav

#ble line in
amixer cset numid=46 1
amixer cset numid=43 1

# Start mua service
# Temp put it here, move to a dedicate script eventually
#
# Must start main service first of all

exit 0
